# docker/php/etc/supervisord.d/cron.conf
[program:cron]
command=crond -n
autostart=true
autorestart=true
